# root/etc/s6-overlay/s6-rc.d/svc-airupnp/run
#!/usr/bin/with-contenv bash
exec /bin/airupnp-linux-x86_64 -l 1000:2000