# root/etc/s6-overlay/s6-rc.d/svc-aircast/run
#!/usr/bin/execlineb -P
exec /bin/aircast-linux-x86_64