# root/etc/s6-overlay/s6-rc.d/svc-aircast/run
#!/usr/bin/with-contenv bash

# Check for VAR - kill, custom, default
if [ "$AIRCAST_VAR" = "kill" ]; then
  echo "AIRCAST_VAR set to kill, killing aircast service launch"
  exit 0
elif [ "$AIRCAST_VAR" != "kill" ]; then
  echo "AIRCAST_VAR present, using custom user aircast launch variables"
  exec /bin/aircast-linux-x86_64 "$AIRCAST_VAR"
else;
  echo "No AIRCAST_VAR present, continuing with standard aircast launch variables"
  exec /bin/aircast-linux-x86_64
fi