# root/etc/services.d/airconnect/run
#!/usr/bin/with-contenv bash

supervisord --nodaemon --configuration /etc/supervisord.conf