# root/etc/s6-overlay/s6-rc.d/init-airconnect/run
#!/usr/bin/execlineb -P
exec /etc/s6-overlay/s6-rc.d/init-airconnect/airconnect-setup.sh
