# root/etc/s6-overlay/s6-rc.d/init-airconnect/run
#!/usr/bin/with-contenv bash
exec /etc/s6-overlay/s6-rc.d/init-airconnect/airconnect-setup.sh
